# rpi/imaging/polypod.service
[Unit]
Description=Polypod Flutter App (kiosk)
After=multi-user.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=POLYPOD_USER_PLACEHOLDER
WorkingDirectory=/home/POLYPOD_USER_PLACEHOLDER

# Wayland/DRM environment
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WLR_BACKENDS=drm
Environment=WLR_DRM_NO_ATOMIC=1
Environment=SEATD_VTBOUND=0

# Launch Flutter app inside Cage (single-window Wayland kiosk compositor)
ExecStart=/usr/bin/cage -s -- /opt/polypod/app/polypod_hw

Restart=on-failure
RestartSec=5

# GPU + input access
SupplementaryGroups=video render input

[Install]
WantedBy=multi-user.target
