# rpi/imaging_new/layer/polypod/rootfs-overlay/etc/systemd/system/app-update.timer
[Unit]
Description=Polypod App Update Timer

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
