# rpi/imaging_new/layer/polypod/rootfs-overlay/etc/systemd/system/app-update.service
[Unit]
Description=Polypod App Update Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
ExecStart=/usr/local/bin/polypod-update.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polypod-update
