# rpi/imaging_new/layer/polypod/rootfs-overlay/home/pi/.config/systemd/user/polypod-bottom-spi.service
[Unit]
Description=Polypod Bottom Window — SPI Fallback (direct DRM)
After=graphical-session.target polypod-top.service
Requires=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment=HOME=/home/pi
EnvironmentFile=/etc/polypod/polypod.conf

ExecCondition=/bin/sh -c '[ "$POLYPOD_SPI_FALLBACK" = "1" ]'
ExecStart=/usr/local/bin/polypod-spi-fallback.sh

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=polypod-bottom-spi

[Install]
WantedBy=default.target
